# deploy/dove-backup.timer
[Unit]
Description=Run Dove Postgres backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
